# meta-rpicustom/recipes-graphics/weston/files/weston.service
[Unit]
Description=Weston Wayland Compositor
RequiresMountsFor=/run
Conflicts=getty@tty1.service xorg.service lightdm.service gdm.service sddm.service
After=systemd-user-sessions.service getty@tty1.service plymouth-quit.service multi-user.target graphical-session-pre.target
Before=graphical.target
Wants=multi-user.target
# Don't require framebuffer device since we're using DRM
# Requires=sys-devices-platform-gpu-graphics-fb0.device

[Service]
Type=simple
ExecStart=/usr/bin/weston --log=/var/log/weston.log --idle-time=0
User=root
Group=root
WorkingDirectory=/home/root
PAMName=login

# Ensure XDG_RUNTIME_DIR is set
Environment=XDG_RUNTIME_DIR=/run/user/0
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_SESSION_CLASS=user
Environment=XDG_SEAT=seat0
Environment=XDG_VTNR=1

# Additional environment variables for stability
Environment=WAYLAND_DEBUG=0
Environment=WESTON_LOG_VERBOSITY=1

# Create runtime directory if it doesn't exist
ExecStartPre=/bin/mkdir -p /run/user/0
ExecStartPre=/bin/chmod 0700 /run/user/0

# Wait for DRM devices to be ready instead of framebuffer
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do [ -e /dev/dri/card0 ] && break || sleep 1; done'

# Wait a bit more before starting to ensure system is ready
ExecStartPre=/bin/sleep 3

StandardOutput=journal
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10s
TimeoutStartSec=60s
Restart=on-failure
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=3

[Install]
WantedBy=graphical.target
Also=graphical.target
